# roles/setup_gitea/templates/rhde_gitops/app_fuxa_scada.container.j2
[Unit]
Description=The FUXA container
After=local-fs.target

[Container]
Image=docker.io/frangoteam/fuxa:latest

[Install]
# Start by default on boot
WantedBy=multi-user.target default.target